# hyperion-elite-bot.service
[Unit]
Description=HYPERION Elite Bot - Headless Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=hyperion
Group=hyperion
WorkingDirectory=/opt/hyperion-elite-bot
ExecStart=/opt/hyperion-elite-bot/venv/bin/python /opt/hyperion-elite-bot/hyperion_headless.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONPATH=/opt/hyperion-elite-bot
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hyperion-elite-bot /var/log/hyperion /tmp
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target